# native/systemd/facephys.service
[Unit]
Description=FacePhys native rPPG monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
SupplementaryGroups=spi dialout
WorkingDirectory=/home/pi/FacePhys-Demo
ExecStart=/home/pi/FacePhys-Demo/native/scripts/run_facephys.sh
Restart=on-failure
RestartSec=3
Nice=0

[Install]
WantedBy=multi-user.target
